# build/supervisord.conf
[supervisord]
nodaemon=true
user=www-data

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true